# roles/software/files/etc/anacrontab
# /etc/anacrontab: configuration file for anacron

# See anacron(8) and anacrontab(5) for details.

SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
HOME=/root
LOGNAME=root
# the maximal random delay added to the base delay of the jobs
RANDOM_DELAY=30
# the jobs will be started during the following hours only
START_HOURS_RANGE=6-23
MAILTO=root@jimbo.wiki.muc

# These replace cron's entries
1	5	cron.daily	nice run-parts --report /etc/cron.daily
7	10	cron.weekly	nice run-parts --report /etc/cron.weekly
@monthly	15	cron.monthly	nice run-parts --report /etc/cron.monthly
